# sorting.socket
[Unit]
Description=sorting socket

[Socket]
ListenStream=/run/sorting.sock

[Install]
WantedBy=sockets.target